# config/self-host/boring-auto-deploy.service.template
[Unit]
Description=Auto deploy boring-ui full-app from GitHub branch
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=oneshot
User=root
Environment=BORING_DEPLOY_BRANCH=plan/self-host-vm-boring
Environment=BORING_DEPLOY_REPO=https://github.com/hachej/boring-ui.git
Environment=BORING_APP_DIR=/opt/boring/full-app
Environment=BORING_IMAGE=boring-full-app:self-host-test
ExecStart=/usr/local/bin/boring-auto-deploy
